# example_configs/usr/lib/systemd/system/reverse_proxy_service.service
[Unit]
Description=reverse proxy service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/computer/workspace/reverse_proxy --c /etc/reverse_proxy/service1.ini
Restart=on-failure

[Install]
WantedBy=multi-user.target
